# lib/systemd/system/distroupdate-automation-autoremove.service
[Unit]
Description=Update Manager weekly autoremove

[Service]
Type=oneshot
ExecStart=/usr/share/distroupdate/automation/autoremove.sh
